# miniconda3/etc/profile.d/conda.csh
setenv CONDA_EXE "/home/ec2-user/miniconda3/bin/conda";
setenv _CONDA_ROOT "/home/ec2-user/miniconda3";
setenv _CONDA_EXE "/home/ec2-user/miniconda3/bin/conda";
setenv CONDA_PYTHON_EXE "/home/ec2-user/miniconda3/bin/python";
source "/home/ec2-user/miniconda3/lib/python3.12/site-packages/conda/shell/etc/profile.d/conda.csh";
